# supervisord.conf
[supervisord]
nodaemon=true

[program:samba]
command=/usr/sbin/samba -i -s /var/lib/samba/private/smb.conf

[program:syslog]
command=/usr/sbin/rsyslogd -n
